# templates/default/systemd.service.erb
[Unit]
Description=<%= @name %>
Wants=network.target
After=network.target

[Service]
Environment=<%= @environment.map {|key, val| %Q{"#{key}=#{val}"} }.join(' ') %>
ExecStart=<%= @command %>
ExecReload=/bin/kill -HUP $MAINPID
User=<%= @service_user %>
LimitNOFILE=65536
Restart=on-failure

[Install]
WantedBy=multi-user.target
